# deploy/cf-excel.timer
[Unit]
Description=Content Factory — excel-конвейер: тик каждые 10 минут

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
