# deploy/cf-mail.timer
[Unit]
Description=Content Factory — почта: проверка каждые 30 минут

[Timer]
OnCalendar=*:5/30
Persistent=true

[Install]
WantedBy=timers.target
