# deploy/cf-scheduler.service
[Unit]
Description=Content Factory — scheduler (исполнение дозревших слотов задач)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/content-factory
Environment=PYTHONPATH=/opt/content-factory/src
ExecStart=/opt/content-factory/.venv/bin/python -m content_factory.orchestrator.scheduler_run
